# etc/ch-sf.conf
# ch-sf.conf
#
# start serialforwarder
#
# Part of cogent-house
#

description "TinyOS SerialForwarder"
author "J. Brusey <james.brusey@gmail.com>"

start on runlevel [2345]

stop on runlevel [016]

respawn

script 
  exec > /tmp/ch-sf.out 2>&1
  . /opt/tinyos-main-read-only/tinyos.sh
  exec java net.tinyos.sf.SerialForwarder -comm serial@/dev/ttyUSB0:telosb -no-output
end script
